# etc/init/sled.conf

description "Sled Server"
author "Ivar Clemens <post@ivarclemens.nl>"

start on runlevel [345]
stop on runlevel [0126]

respawn

expect fork

limit memlock 1048576
exec sled-server-1.1

